# packaging/systemd/hindsight-worker.service
[Unit]
Description=Hindsight Worker (async task processor)
After=network.target hindsight-api.service
Wants=hindsight-api.service

[Service]
Type=simple
User=hindsight
Group=hindsight
WorkingDirectory=/var/lib/hindsight
Environment="HOME=/var/lib/hindsight"
Environment="PATH=/opt/hindsight/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/opt/hindsight/venv/bin/hindsight-worker
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hindsight-worker
TimeoutStopSec=30s

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/hindsight

[Install]
WantedBy=multi-user.target
